# contrib/ferrosonicd.service
[Unit]
Description=Ferrosonic music daemon
Documentation=https://github.com/Jamie098/ferrosonic
After=graphical-session.target pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/ferrosonic --daemon
Restart=on-failure
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=default.target
